# yacey-graphite.timer
[Unit]
Description=Send YaCy stats to Graphite

[Timer]
OnCalendar=*:*:0/10
Unit=yacy-graphite.service

[Install]
WantedBy=yacy-graphite.service
