# yacey-graphite.timer
[Unit]
Description=Send YaCy stats to Graphite

[Timer]
OnCalendar=10s
Unit=yacey-graphite.service

[Install]
WantedBy=yacey-graphite.service
